# ffmpeg-rtsp-hls.path
[Unit]
Description=Path-based activation unit for ffmpeg-rtsp-hls.service
After=network.target

[Path]
PathExists=/path/camera-identity/active

[Install]
WantedBy=multi-user.target

